# openclaw-management/scripts/systemd/openclaw-management.service
[Unit]
Description=OpenClaw 管理设置页 (量产 8080 / 开发可改 8000)
Documentation=file:///opt/amyclaw/openclaw-management/README.md
After=network-online.target openclaw-gateway.service
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
# 量产默认 WorkingDirectory=/opt/amyclaw/openclaw-management；开发机可用 install-autostart.sh 传入其它仓库根路径以替换单元内 /opt/amyclaw
WorkingDirectory=/opt/amyclaw/openclaw-management
ExecStart=/usr/bin/env node server.js
Restart=on-failure
RestartSec=5
# 与 OpenClaw Gateway 共用状态目录；必须与 gateway 一致
Environment=OPENCLAW_STATE_DIR=/root/.openclaw
# 量产默认 8080；开发热加载调试可用 MANAGEMENT_PORT=8000
Environment=MANAGEMENT_PORT=8080

[Install]
WantedBy=multi-user.target
