# openclaw-management/scripts/systemd/openclaw-management.service
[Unit]
Description=OpenClaw 管理设置页 (量产 8080 / 开发可改 8000)
Documentation=file:///mnt/disk/amyclaw/openclaw-management/README.md
After=network-online.target openclaw-gateway.service
Wants=network-online.target

[Service]
Type=simple
# 开发：指向仓库路径；量产：使用 scripts/deploy-to-opt.sh 部署到 /opt/openclaw-management 后本单元会被覆盖为 WorkingDirectory=/opt/openclaw-management
WorkingDirectory=/mnt/disk/amyclaw/openclaw-management
ExecStart=/usr/bin/env node server.js
Restart=on-failure
RestartSec=5
# 与 OpenClaw Gateway 共用状态目录；必须与 gateway 一致
Environment=OPENCLAW_STATE_DIR=/root/.openclaw
# 量产默认 8080；开发热加载调试可用 MANAGEMENT_PORT=8000
Environment=MANAGEMENT_PORT=8080

[Install]
WantedBy=multi-user.target
